# config/backend_sz_app.conf
description "Gunicorn application server running sz backend app"
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid root
setgid www-data

env PATH=/opt/sz_app/app/env/bin
chdir /opt/sz_app
exec gunicorn --workers 3 --bind unix:sz.sock -m 007 wsgi:app
